# greengrass/ggc_groups_ste/device/sample_zero/ggdBTNPressYellow/systemd.service.j2
[Unit]
Description={{ device_name }} Daemon

[Service]
Type=idle
ExecStart=/gg_devices/{{ device_name }}/venv/bin/python /gg_devices/{{ device_name }}/ggdBTNPressYellow.py
KillSignal=SIGQUIT

[Install]
WantedBy=multi-user.target
